# deploy/barbercrm.service
[Unit]
Description=BarberCRM Django Backend (Gunicorn)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/barbercrm/backend
ExecStart=/var/www/barbercrm/venv/bin/gunicorn \
    config.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --timeout 120 \
    --access-logfile /var/log/barbercrm/gunicorn-access.log \
    --error-logfile /var/log/barbercrm/gunicorn-error.log \
    --capture-output \
    --log-level info
Restart=always
RestartSec=5
EnvironmentFile=/var/www/barbercrm/backend/.env

[Install]
WantedBy=multi-user.target
